# modbus_tcp_test.conf
[program:modbus_tcp_test]

command=python /home/weivi/program/modbus_tcp_test/run.py

autostart=true

autorestart=ture

stdout_logfile=/home/weivi/log/modbus_tcp_test/test.log